# templates/minecraft-backup.service.j2
[Unit]
Description=Minecraft World Backup

[Service]
Type=simple
ExecStart={{ minecraft_dir }}/tools/backup.sh
