# files/minecraft-backup.service
[Unit]
Description=Minecraft World Backup

[Service]
Type=simple
ExecStart=/opt/minecraft/tools/backup.sh
